# templates/etc/systemd/system/docker.service.d/docker.conf.j2
# {{ ansible_managed }}

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd
Environment=AWS_SHARED_CREDENTIALS_FILE=/root/.aws/credentials

